# fruitwars.conf
description "Fruit Wars!"
author "John Shaw <jshaw@shawnation.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /opt/fruitwars
